# com.vaultassistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.vaultassistant</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/env</string>
		<string>uv</string>
		<string>run</string>
		<string>uvicorn</string>
		<string>main:app</string>
		<string>--host</string>
		<string>0.0.0.0</string>
		<string>--port</string>
		<string>8765</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Volumes/1TB/Repos/vault-assistant</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>~/.vault-assistant/logs/stdout.log</string>
	<key>StandardErrorPath</key>
	<string>~/.vault-assistant/logs/stderr.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
</dict>
</plist>
